# .config/systemd/user/xsession.target
[Unit]
Description=User X session
After=default.target
AllowIsolate=yes
